# indiainnovates-automation.service
[Unit]
Description=India Innovates 2026 Evaluator - Instance %i
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/indiainnovates-automation
Environment="PATH=/root/indiainnovates-automation/venv/bin"
ExecStart=/root/indiainnovates-automation/venv/bin/streamlit run ii2026_evaluator.py \
    --server.port=850%i \
    --server.address=127.0.0.1 \
    --server.headless=true \
    --server.enableCORS=false \
    --server.enableXsrfProtection=true
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
